# startups.conf
[program:startups]
# 启动命令（建议使用虚拟环境的绝对路径）
command=uvicorn server.main:app --host :: --port 80

# 项目根目录
directory=/root/StartUp

# 运行服务的用户（不建议用 root）
user=root

# 自动启动和崩溃自动重启
autostart=true
autorestart=true

# 重启延迟和重试次数
startsecs=5
startretries=3

# 日志配置（非常重要，用于排查错误）
stdout_logfile=/var/log/my_game.out.log
stderr_logfile=/var/log/my_game.err.log
stdout_logfile_maxbytes=1GB
stdout_logfile_backups=5